# vagrant/debian/jessie/supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
startretries=100000000
autorestart=true
